# urlswithhappybase/urlparser/test.cron
30 12 */15 * *  ./home/ubuntu/pythonlib/urlswithhappybase/urlparser/main.py
